# supervisord.conf
[supervisord]
nodaemon=true
loglevel=debug

[program:rpicontroller]
command=/usr/src/rpi-controller/bin/controller-server
redirect_stderr=true
