# supervisord.conf
[supervisord]
nodaemon=true
loglevel=debug

[program:rpicontroller]
command=/controller-server
redirect_stderr=true
